# rc.local
#!/bin/sh
cd /home/bwb0de/.scripts
./wpchange&
